# .config/supervisord.conf
[program:gunicorn]
command = gunicorn config.wsgi.prod -c ../.config/gunicorn_cfg.py
directory = /srv/bookcheck/app

[program:nginx]
command = nginx -g 'daemon off;'